# config/_zaborona_v2/etc/systemd/system/openvpn-server@.service.d/override.conf
# LimitNPROC does not work properly in container for some reason.

[Service]
LimitNPROC=infinity 
LimitNOFILE=65535
