# config/_zaborona_v2/etc/systemd/system/zaborona-update.timer
[Unit]
Description=Run zaborona-vpn update every 6h

[Timer]
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target